# crontab
0 4 * * * /usr/local/bin/update_playlist.sh >> /var/log/cron.log 2>&1
